# chess-server.service
[Unit]
Description=Chinese Chess Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/peterc20/opencode/chinese-chcess
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/home/peterc20/.nvm/versions/node/v24.13.1/bin/node /home/peterc20/opencode/chinese-chcess/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
