# chess-frontend.service
[Unit]
Description=Chinese Chess Frontend
After=network.target chess-server.service

[Service]
Type=simple
WorkingDirectory=/home/peterc20/opencode/chinese-chcess/frontend
Environment=NODE_ENV=production
Environment=PATH=/home/peterc20/.nvm/versions/node/v24.13.1/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStartPre=/home/peterc20/.nvm/versions/node/v24.13.1/bin/node /home/peterc20/opencode/chinese-chcess/frontend/node_modules/vite/bin/vite.js build
ExecStart=/home/peterc20/.nvm/versions/node/v24.13.1/bin/node /home/peterc20/opencode/chinese-chcess/frontend/node_modules/vite/bin/vite.js preview --host 0.0.0.0 --port 5173 --strictPort
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
